# supervisor/conf.d/auth-canary.conf
; ============================================================
; JonOps — Claude OAuth Auth Canary
; ============================================================
; Managed by supervisord. Auto-starts on container boot,
; auto-restarts on crash. Runs claude --print -p ok daily at
; 03:00 WITA. Alerts Telegram on non-zero exit so a broken
; refresh chain is caught before operator-facing slots fail.
; ============================================================

[program:auth-canary]
command=bash /home/agent/project/claude-auth-canary.sh
directory=/home/agent/project
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stdout_logfile=/home/agent/project/logs/auth-canary-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_redirect=true
environment=HOME="/home/agent",PATH="/home/agent/.local/bin:/usr/local/bin:/usr/bin:/bin"
