# supervisor/conf.d/scheduler.conf
; ============================================================
; JonOps — Scheduler Process
; ============================================================
; Managed by supervisord. Auto-starts on container boot,
; auto-restarts on crash.
; ============================================================

[program:scheduler]
command=python3 /home/agent/project/scheduler.py --force
directory=/home/agent/project
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stdout_logfile=/home/agent/project/logs/scheduler-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_redirect=true
environment=HOME="/home/agent",PATH="/home/agent/.local/bin:/usr/local/bin:/usr/bin:/bin"
