# supervisor/conf.d/telegram-daemon.conf
; ============================================================
; JonOps — Telegram Chat Daemon
; ============================================================
; Managed by supervisord. Auto-starts on container boot,
; auto-restarts on crash. Long-polls Telegram and routes
; messages to the Claude Agent SDK.
; ============================================================

[program:telegram-daemon]
command=python3 /home/agent/project/telegram-daemon.py
directory=/home/agent/project
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stdout_logfile=/home/agent/project/logs/telegram-daemon-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_redirect=true
environment=HOME="/home/agent",PATH="/home/agent/.local/bin:/usr/local/bin:/usr/bin:/bin"
